# x86_64-ssr+_v2s/files/etc/uci-defaults/99-custom
#!/bin/sh

uci set system.@system[0].zonename=Asia/Shanghai
uci commit system
