# x86_64-25.12.5-daede/files/etc/uci-defaults/99-custom
#!/bin/sh

uci set system.@system[0].zonename=Asia/Shanghai
uci commit system

sed -i -e '/daede/d' /etc/apk/repositories.d/distfeeds.list 2>/dev/null || true

exit 0
